# etc/udev/rules.d/99-usb-keyboards.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="04d9", ATTRS{idProduct}=="0141", RUN+="/home/glider/bin/kbd-udev", OWNER="glider"
